# deploy/com.fireguard.serve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  FGDC live SERVE LaunchAgent (Ship A). Fires 4×/day at ~06:15, 07:15, 18:15, 19:15 LOCAL (CEST), plus at
  load/login. Each slot is timed ~30 min after a VIIRS pass settles in FIRMS NRT (≈3 h latency), so every run
  folds the newest active-fire pass into today's dist_to_fire before predicting t+1 — the "use the freshest
  info as soon as we can" cadence. Three VIIRS birds (S-NPP + NOAA-20 + NOAA-21) → 6 daily passes in two clusters:
  the night cluster (~00:00-01:30 UTC → picked up ~04/05 UTC = 06/07 local) and the afternoon cluster
  (~12:00-13:30 UTC → ~16/17 UTC = 18/19 local). NOAA-21 LEADS the group so its passes settle earliest — the 4 slots
  (timed to the LATEST pass of each cluster, S-NPP) already fold in the earlier NOAA-20/21 passes, so the 3rd bird
  needs NO new slots. The first three runs land PRELIMINARY; the 19:15 run (>17 UTC settle) is the FINAL for t+1.
  run_serve.sh predicts + pushes to the HF Dataset (the Space reads it), dedups RunAtLoad vs a slot (MIN_GAP_MIN,
  sub-hour so the 1 h-apart pairs both fire) and catches up once on wake — a missed slot is not backfilled (serve
  always predicts the latest). Ephemeral: writes nothing to the cube. User-session agent (needs the
  repo/.venv/.env/network + the gold cube on this machine).

  Install:   cp deploy/com.fireguard.serve.plist ~/Library/LaunchAgents/
             launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.fireguard.serve.plist
  Status:    launchctl print gui/$(id -u)/com.fireguard.serve
  Run now:   launchctl kickstart -k gui/$(id -u)/com.fireguard.serve
  Uninstall: launchctl bootout gui/$(id -u)/com.fireguard.serve && rm ~/Library/LaunchAgents/com.fireguard.serve.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.fireguard.serve</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/vladimir/ml-wildfire-prediction/scripts/run_serve.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/vladimir/ml-wildfire-prediction</string>
  <key>RunAtLoad</key>
  <true/>
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>6</integer>
      <key>Minute</key>
      <integer>15</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>7</integer>
      <key>Minute</key>
      <integer>15</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>18</integer>
      <key>Minute</key>
      <integer>15</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>19</integer>
      <key>Minute</key>
      <integer>15</integer>
    </dict>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/bin:/bin:/usr/sbin:/sbin:/usr/local/bin:/opt/homebrew/bin</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/vladimir/ml-wildfire-prediction/data/serving_store/serve.launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/vladimir/ml-wildfire-prediction/data/serving_store/serve.launchd.err.log</string>
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
